# ryzenadj.service
[Unit]
Description=Ryzenadj

[Service]
ExecStart=/usr/local/sbin/ryzenadj.sh


[Install]
WantedBy=multi-user.target
